# scripts/systemd/nanoclaw-audit-v1-bugs.service
[Unit]
Description=NanoClaw v1-bug auditor (daily)
After=network.target nanoclaw.service
# Doesn't Require=nanoclaw.service — the auditor opens taskflow.db
# read-only, so it can run even if the host service is restarting.

[Service]
Type=oneshot
User=nanoclaw
Group=nanoclaw
WorkingDirectory=/root/nanoclaw
ExecStart=/root/nanoclaw/node_modules/.bin/tsx /root/nanoclaw/scripts/audit-v1-bugs-daily.ts
StandardOutput=append:/root/nanoclaw/logs/audit-v1-bugs.log
StandardError=append:/root/nanoclaw/logs/audit-v1-bugs.error.log
Environment=PATH=/usr/local/bin:/usr/bin:/bin
# Hardening. `ProtectHome=read-only` is intentional (not `true`): the
# install lives under /root/nanoclaw, and `ProtectHome=true` would make
# the whole /root tree appear empty so ExecStart can't even resolve.
# read-only still blocks accidental writes into other users' homes;
# ReadWritePaths then carves out the audit-output + log directories.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/root/nanoclaw/data/audit /root/nanoclaw/logs
PrivateTmp=true
# Operational fuses. The detector is sub-second on the pilot DB, but
# task_history grows with traffic and the self-joins are exactly the
# class of audit query that should have a bound — a runaway audit must
# not hold a shared lock long enough to back up host writers
# (taskflow.db uses journal_mode=DELETE; long shared locks block
# committers).
TimeoutSec=300
MemoryMax=512M

[Install]
WantedBy=multi-user.target
